# config/etc/cron.d/audiencias
* * * * * python3 /var/labhacker/audiencias/manage.py get_camara_webservice
* * * * * python3 /var/labhacker/audiencias/manage.py prune_rooms
* * * * * python3 /var/labhacker/audiencias/manage.py prune_presences
